# ansible/roles/configure/system76/charge-thresholds/files/charge-thresholds.service
[Unit]
Description=Set the charge threshold at startup.
After=default.target

[Service]
Type=simple
ExecStart=system76-power charge-thresholds --profile max_lifespan

[Install]
WantedBy=default.target
